# examples/crontab
25 6 * * * /usr/sbin/scrape-credentials
